# base/soc/qm10xd/linux/prebuilts/esl/bsp/initramfs.list

#####################
# /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/
# Last modified: 1730892574.6185002530

dir /var 755 123 127
dir /usr 755 123 127
dir /usr/sbin 755 123 127
slink /usr/sbin/i2ctransfer ../../bin/busybox 777 123 127
slink /usr/sbin/fdformat ../../bin/busybox 777 123 127
slink /usr/sbin/fbset ../../bin/busybox 777 123 127
slink /usr/sbin/chroot ../../bin/busybox 777 123 127
slink /usr/sbin/dnsd ../../bin/busybox 777 123 127
slink /usr/sbin/crond ../../bin/busybox 777 123 127
slink /usr/sbin/rdate ../../bin/busybox 777 123 127
slink /usr/sbin/i2cdetect ../../bin/busybox 777 123 127
slink /usr/sbin/partprobe ../../bin/busybox 777 123 127
slink /usr/sbin/killall5 ../../bin/busybox 777 123 127
slink /usr/sbin/ether-wake ../../bin/busybox 777 123 127
slink /usr/sbin/arping ../../bin/busybox 777 123 127
slink /usr/sbin/mim ../../bin/busybox 777 123 127
slink /usr/sbin/ubirename ../../bin/busybox 777 123 127
slink /usr/sbin/telnetd ../../bin/busybox 777 123 127
slink /usr/sbin/i2cget ../../bin/busybox 777 123 127
slink /usr/sbin/inetd ../../bin/busybox 777 123 127
slink /usr/sbin/deluser ../../bin/busybox 777 123 127
slink /usr/sbin/loadfont ../../bin/busybox 777 123 127
slink /usr/sbin/readprofile ../../bin/busybox 777 123 127
slink /usr/sbin/adduser ../../bin/busybox 777 123 127
slink /usr/sbin/fsfreeze ../../bin/busybox 777 123 127
slink /usr/sbin/seedrng ../../bin/busybox 777 123 127
slink /usr/sbin/i2cset ../../bin/busybox 777 123 127
slink /usr/sbin/i2cdump ../../bin/busybox 777 123 127
slink /usr/sbin/delgroup ../../bin/busybox 777 123 127
slink /usr/sbin/setlogcons ../../bin/busybox 777 123 127
slink /usr/sbin/addgroup ../../bin/busybox 777 123 127
slink /usr/sbin/nologin ../../bin/busybox 777 123 127
dir /usr/share 755 123 127
dir /usr/share/udhcpc 755 123 127
file /usr/share/udhcpc/default.script /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/usr/share/udhcpc/default.script 755 123 127
dir /usr/lib 755 123 127
dir /usr/bin 755 123 127
slink /usr/bin/paste ../../bin/busybox 777 123 127
slink /usr/bin/microcom ../../bin/busybox 777 123 127
slink /usr/bin/lzopcat ../../bin/busybox 777 123 127
slink /usr/bin/hostid ../../bin/busybox 777 123 127
slink /usr/bin/tail ../../bin/busybox 777 123 127
slink /usr/bin/svok ../../bin/busybox 777 123 127
slink /usr/bin/passwd ../../bin/busybox 777 123 127
slink /usr/bin/uniq ../../bin/busybox 777 123 127
slink /usr/bin/cmp ../../bin/busybox 777 123 127
slink /usr/bin/head ../../bin/busybox 777 123 127
slink /usr/bin/chrt ../../bin/busybox 777 123 127
slink /usr/bin/crontab ../../bin/busybox 777 123 127
slink /usr/bin/unlzop ../../bin/busybox 777 123 127
slink /usr/bin/reset ../../bin/busybox 777 123 127
slink /usr/bin/unlink ../../bin/busybox 777 123 127
slink /usr/bin/[ ../../bin/busybox 777 123 127
slink /usr/bin/ascii ../../bin/busybox 777 123 127
slink /usr/bin/fold ../../bin/busybox 777 123 127
slink /usr/bin/expr ../../bin/busybox 777 123 127
slink /usr/bin/xxd ../../bin/busybox 777 123 127
slink /usr/bin/lsscsi ../../bin/busybox 777 123 127
slink /usr/bin/xzcat ../../bin/busybox 777 123 127
slink /usr/bin/id ../../bin/busybox 777 123 127
slink /usr/bin/bzcat ../../bin/busybox 777 123 127
slink /usr/bin/sha512sum ../../bin/busybox 777 123 127
slink /usr/bin/sha256sum ../../bin/busybox 777 123 127
slink /usr/bin/less ../../bin/busybox 777 123 127
slink /usr/bin/free ../../bin/busybox 777 123 127
slink /usr/bin/ipcs ../../bin/busybox 777 123 127
slink /usr/bin/ipcrm ../../bin/busybox 777 123 127
slink /usr/bin/tree ../../bin/busybox 777 123 127
slink /usr/bin/pkill ../../bin/busybox 777 123 127
slink /usr/bin/xargs ../../bin/busybox 777 123 127
slink /usr/bin/printf ../../bin/busybox 777 123 127
slink /usr/bin/awk ../../bin/busybox 777 123 127
slink /usr/bin/tr ../../bin/busybox 777 123 127
slink /usr/bin/hexedit ../../bin/busybox 777 123 127
slink /usr/bin/tftp ../../bin/busybox 777 123 127
slink /usr/bin/uptime ../../bin/busybox 777 123 127
slink /usr/bin/mkpasswd ../../bin/busybox 777 123 127
slink /usr/bin/renice ../../bin/busybox 777 123 127
slink /usr/bin/du ../../bin/busybox 777 123 127
slink /usr/bin/wget ../../bin/busybox 777 123 127
slink /usr/bin/eject ../../bin/busybox 777 123 127
slink /usr/bin/unix2dos ../../bin/busybox 777 123 127
slink /usr/bin/env ../../bin/busybox 777 123 127
slink /usr/bin/yes ../../bin/busybox 777 123 127
slink /usr/bin/ts ../../bin/busybox 777 123 127
slink /usr/bin/lzcat ../../bin/busybox 777 123 127
slink /usr/bin/lsusb ../../bin/busybox 777 123 127
slink /usr/bin/cut ../../bin/busybox 777 123 127
slink /usr/bin/setkeycodes ../../bin/busybox 777 123 127
slink /usr/bin/readlink ../../bin/busybox 777 123 127
slink /usr/bin/telnet ../../bin/busybox 777 123 127
slink /usr/bin/w ../../bin/busybox 777 123 127
slink /usr/bin/setsid ../../bin/busybox 777 123 127
slink /usr/bin/dos2unix ../../bin/busybox 777 123 127
slink /usr/bin/nohup ../../bin/busybox 777 123 127
slink /usr/bin/whoami ../../bin/busybox 777 123 127
slink /usr/bin/traceroute ../../bin/busybox 777 123 127
slink /usr/bin/basename ../../bin/busybox 777 123 127
slink /usr/bin/time ../../bin/busybox 777 123 127
slink /usr/bin/tsort ../../bin/busybox 777 123 127
slink /usr/bin/mkfifo ../../bin/busybox 777 123 127
slink /usr/bin/who ../../bin/busybox 777 123 127
slink /usr/bin/setfattr ../../bin/busybox 777 123 127
slink /usr/bin/fallocate ../../bin/busybox 777 123 127
slink /usr/bin/find ../../bin/busybox 777 123 127
slink /usr/bin/tee ../../bin/busybox 777 123 127
slink /usr/bin/fuser ../../bin/busybox 777 123 127
slink /usr/bin/patch ../../bin/busybox 777 123 127
slink /usr/bin/killall ../../bin/busybox 777 123 127
slink /usr/bin/unzip ../../bin/busybox 777 123 127
slink /usr/bin/svc ../../bin/busybox 777 123 127
slink /usr/bin/realpath ../../bin/busybox 777 123 127
slink /usr/bin/uuencode ../../bin/busybox 777 123 127
slink /usr/bin/install ../../bin/busybox 777 123 127
slink /usr/bin/dc ../../bin/busybox 777 123 127
slink /usr/bin/pmap ../../bin/busybox 777 123 127
slink /usr/bin/xz ../../bin/busybox 777 123 127
slink /usr/bin/which ../../bin/busybox 777 123 127
slink /usr/bin/hexdump ../../bin/busybox 777 123 127
slink /usr/bin/sha3sum ../../bin/busybox 777 123 127
slink /usr/bin/bc ../../bin/busybox 777 123 127
slink /usr/bin/factor ../../bin/busybox 777 123 127
slink /usr/bin/deallocvt ../../bin/busybox 777 123 127
slink /usr/bin/unxz ../../bin/busybox 777 123 127
slink /usr/bin/md5sum ../../bin/busybox 777 123 127
slink /usr/bin/last ../../bin/busybox 777 123 127
slink /usr/bin/openvt ../../bin/busybox 777 123 127
slink /usr/bin/chvt ../../bin/busybox 777 123 127
slink /usr/bin/pgrep ../../bin/busybox 777 123 127
slink /usr/bin/clear ../../bin/busybox 777 123 127
slink /usr/bin/lsof ../../bin/busybox 777 123 127
slink /usr/bin/strings ../../bin/busybox 777 123 127
slink /usr/bin/diff ../../bin/busybox 777 123 127
slink /usr/bin/vlock ../../bin/busybox 777 123 127
slink /usr/bin/flock ../../bin/busybox 777 123 127
slink /usr/bin/truncate ../../bin/busybox 777 123 127
slink /usr/bin/nl ../../bin/busybox 777 123 127
slink /usr/bin/wc ../../bin/busybox 777 123 127
slink /usr/bin/logname ../../bin/busybox 777 123 127
slink /usr/bin/logger ../../bin/busybox 777 123 127
slink /usr/bin/sha1sum ../../bin/busybox 777 123 127
slink /usr/bin/resize ../../bin/busybox 777 123 127
slink /usr/bin/dirname ../../bin/busybox 777 123 127
slink /usr/bin/bunzip2 ../../bin/busybox 777 123 127
slink /usr/bin/unlzma ../../bin/busybox 777 123 127
slink /usr/bin/uudecode ../../bin/busybox 777 123 127
slink /usr/bin/lspci ../../bin/busybox 777 123 127
slink /usr/bin/crc32 ../../bin/busybox 777 123 127
slink /usr/bin/cksum ../../bin/busybox 777 123 127
slink /usr/bin/mesg ../../bin/busybox 777 123 127
slink /usr/bin/od ../../bin/busybox 777 123 127
slink /usr/bin/nproc ../../bin/busybox 777 123 127
slink /usr/bin/tty ../../bin/busybox 777 123 127
slink /usr/bin/seq ../../bin/busybox 777 123 127
slink /usr/bin/nslookup ../../bin/busybox 777 123 127
slink /usr/bin/shred ../../bin/busybox 777 123 127
slink /usr/bin/top ../../bin/busybox 777 123 127
slink /usr/bin/lzma ../../bin/busybox 777 123 127
slink /usr/bin/sort ../../bin/busybox 777 123 127
slink /usr/bin/[[ ../../bin/busybox 777 123 127
slink /usr/bin/test ../../bin/busybox 777 123 127
slink /usr/bin/ar ../../bin/busybox 777 123 127
dir /home 755 123 127
dir /tmp 755 123 127
dir /ko 755 123 127
dir /sbin 755 123 127
slink /sbin/losetup ../bin/busybox 777 123 127
slink /sbin/poweroff ../bin/busybox 777 123 127
slink /sbin/setconsole ../bin/busybox 777 123 127
slink /sbin/uevent ../bin/busybox 777 123 127
slink /sbin/fsck ../bin/busybox 777 123 127
slink /sbin/sysctl ../bin/busybox 777 123 127
slink /sbin/klogd ../bin/busybox 777 123 127
slink /sbin/iptunnel ../bin/busybox 777 123 127
slink /sbin/tc ../bin/busybox 777 123 127
slink /sbin/ip ../bin/busybox 777 123 127
slink /sbin/pivot_root ../bin/busybox 777 123 127
slink /sbin/switch_root ../bin/busybox 777 123 127
slink /sbin/freeramdisk ../bin/busybox 777 123 127
slink /sbin/mkdosfs ../bin/busybox 777 123 127
slink /sbin/swapoff ../bin/busybox 777 123 127
slink /sbin/runlevel ../bin/busybox 777 123 127
slink /sbin/blkid ../bin/busybox 777 123 127
slink /sbin/syslogd ../bin/busybox 777 123 127
slink /sbin/watchdog ../bin/busybox 777 123 127
slink /sbin/iprule ../bin/busybox 777 123 127
slink /sbin/fdisk ../bin/busybox 777 123 127
slink /sbin/halt ../bin/busybox 777 123 127
slink /sbin/mkfs.ext2 ../bin/busybox 777 123 127
slink /sbin/mke2fs ../bin/busybox 777 123 127
slink /sbin/devmem ../bin/busybox 777 123 127
slink /sbin/init ../bin/busybox 777 123 127
slink /sbin/makedevs ../bin/busybox 777 123 127
slink /sbin/ipaddr ../bin/busybox 777 123 127
slink /sbin/sulogin ../bin/busybox 777 123 127
slink /sbin/rmmod ../bin/busybox 777 123 127
slink /sbin/ipneigh ../bin/busybox 777 123 127
slink /sbin/nameif ../bin/busybox 777 123 127
slink /sbin/udhcpc ../bin/busybox 777 123 127
slink /sbin/insmod ../bin/busybox 777 123 127
slink /sbin/mkswap ../bin/busybox 777 123 127
slink /sbin/swapon ../bin/busybox 777 123 127
slink /sbin/arp ../bin/busybox 777 123 127
slink /sbin/vconfig ../bin/busybox 777 123 127
slink /sbin/reboot ../bin/busybox 777 123 127
slink /sbin/start-stop-daemon ../bin/busybox 777 123 127
slink /sbin/modprobe ../bin/busybox 777 123 127
slink /sbin/route ../bin/busybox 777 123 127
slink /sbin/ifup ../bin/busybox 777 123 127
slink /sbin/getty ../bin/busybox 777 123 127
slink /sbin/iproute ../bin/busybox 777 123 127
slink /sbin/hwclock ../bin/busybox 777 123 127
slink /sbin/hdparm ../bin/busybox 777 123 127
slink /sbin/lsmod ../bin/busybox 777 123 127
slink /sbin/iplink ../bin/busybox 777 123 127
slink /sbin/mdev ../bin/busybox 777 123 127
slink /sbin/loadkmap ../bin/busybox 777 123 127
slink /sbin/ifdown ../bin/busybox 777 123 127
slink /sbin/fstrim ../bin/busybox 777 123 127
slink /sbin/ifconfig ../bin/busybox 777 123 127
slink /sbin/run-init ../bin/busybox 777 123 127
slink /sbin/mkfs.vfat ../bin/busybox 777 123 127
dir /root 755 123 127
file /root/.profile /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/root/.profile 644 123 127
slink /init bin/busybox 777 123 127
dir /dev 755 123 127
nod /dev/null 600 0 0 c 1 3
nod /dev/console 600 0 0 c 5 1
dir /sys 755 123 127
dir /lib 755 123 127
file /lib/libuClibc-1.0.50.so /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/lib/libuClibc-1.0.50.so 755 123 127
slink /lib/libc.so.1 libuClibc-1.0.50.so 777 123 127
slink /lib/libc.so.0 libc.so.1 777 123 127
slink /lib/ld-uClibc.so.0 ld-uClibc.so.1 777 123 127
file /lib/ld-uClibc-1.0.50.so /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/lib/ld-uClibc-1.0.50.so 755 123 127
slink /lib/ld-uClibc.so.1 ld-uClibc-1.0.50.so 777 123 127
dir /bin 755 123 127
slink /bin/fgrep busybox 777 123 127
slink /bin/linux32 busybox 777 123 127
slink /bin/vi busybox 777 123 127
slink /bin/printenv busybox 777 123 127
slink /bin/grep busybox 777 123 127
slink /bin/linux64 busybox 777 123 127
slink /bin/mknod busybox 777 123 127
slink /bin/netstat busybox 777 123 127
slink /bin/tar busybox 777 123 127
slink /bin/nuke busybox 777 123 127
slink /bin/mountpoint busybox 777 123 127
slink /bin/ps busybox 777 123 127
slink /bin/mkdir busybox 777 123 127
slink /bin/rmdir busybox 777 123 127
slink /bin/fdflush busybox 777 123 127
slink /bin/setpriv busybox 777 123 127
slink /bin/hostname busybox 777 123 127
slink /bin/usleep busybox 777 123 127
slink /bin/ln busybox 777 123 127
slink /bin/base32 busybox 777 123 127
slink /bin/setarch busybox 777 123 127
slink /bin/gunzip busybox 777 123 127
slink /bin/arch busybox 777 123 127
slink /bin/dd busybox 777 123 127
slink /bin/mpstat busybox 777 123 127
file /bin/memtester /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/bin/memtester 755 123 127
slink /bin/lsattr busybox 777 123 127
slink /bin/sed busybox 777 123 127
slink /bin/chgrp busybox 777 123 127
slink /bin/mktemp busybox 777 123 127
slink /bin/cp busybox 777 123 127
slink /bin/ping busybox 777 123 127
slink /bin/pipe_progress busybox 777 123 127
file /bin/udevadm /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/bin/udevadm 755 123 127
slink /bin/date busybox 777 123 127
slink /bin/egrep busybox 777 123 127
slink /bin/zcat busybox 777 123 127
slink /bin/gzip busybox 777 123 127
slink /bin/su busybox 777 123 127
slink /bin/more busybox 777 123 127
slink /bin/chmod busybox 777 123 127
slink /bin/umount busybox 777 123 127
slink /bin/cpio busybox 777 123 127
slink /bin/echo busybox 777 123 127
slink /bin/sleep busybox 777 123 127
slink /bin/df busybox 777 123 127
slink /bin/cat busybox 777 123 127
file /bin/busybox /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/bin/busybox 755 123 127
slink /bin/stty busybox 777 123 127
slink /bin/false busybox 777 123 127
slink /bin/dnsdomainname busybox 777 123 127
file /bin/autologin /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/bin/autologin 755 123 127
slink /bin/base64 busybox 777 123 127
slink /bin/resume busybox 777 123 127
slink /bin/chattr busybox 777 123 127
slink /bin/kill busybox 777 123 127
slink /bin/chown busybox 777 123 127
slink /bin/dmesg busybox 777 123 127
slink /bin/ls busybox 777 123 127
slink /bin/run-parts busybox 777 123 127
slink /bin/sync busybox 777 123 127
file /bin/udevd /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/bin/udevd 755 123 127
slink /bin/ash busybox 777 123 127
slink /bin/pwd busybox 777 123 127
slink /bin/touch busybox 777 123 127
slink /bin/nice busybox 777 123 127
slink /bin/mount busybox 777 123 127
slink /bin/pidof busybox 777 123 127
slink /bin/mv busybox 777 123 127
slink /bin/iostat busybox 777 123 127
slink /bin/setserial busybox 777 123 127
slink /bin/login busybox 777 123 127
slink /bin/mt busybox 777 123 127
slink /bin/sh busybox 777 123 127
file /bin/perfmon /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/bin/perfmon 755 123 127
slink /bin/true busybox 777 123 127
slink /bin/dumpkmap busybox 777 123 127
slink /bin/link busybox 777 123 127
slink /bin/getopt busybox 777 123 127
slink /bin/uname busybox 777 123 127
slink /bin/rm busybox 777 123 127
slink /bin/watch busybox 777 123 127
slink /linuxrc bin/busybox 777 123 127
dir /nfs 755 123 127
dir /debug 755 123 127
dir /etc 755 123 127
file /etc/passwd /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/passwd 644 123 127
file /etc/services /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/services 644 123 127
file /etc/group /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/group 644 123 127
dir /etc/init.d 755 123 127
file /etc/init.d/rcS /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/init.d/rcS 755 123 127
file /etc/init.d/Z02mdev /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/init.d/Z02mdev 755 123 127
file /etc/init.d/S01udev /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/init.d/S01udev 755 123 127
file /etc/init.d/S00devs /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/init.d/S00devs 755 123 127
file /etc/fstab /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/fstab 744 123 127
file /etc/profile /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/profile 744 123 127
file /etc/protocols /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/protocols 644 123 127
file /etc/passwd- /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/passwd- 600 123 127
file /etc/inittab /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/inittab 744 123 127
dir /etc/udev 755 123 127
file /etc/udev/usbdev-hotplug.sh /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/usbdev-hotplug.sh 755 123 127
file /etc/udev/disk-hotplug.sh /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/disk-hotplug.sh 755 123 127
file /etc/udev/udev.conf /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/udev.conf 644 123 127
file /etc/udev/firmware.sh /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/firmware.sh 755 123 127
dir /etc/udev/rules.d 755 123 127
file /etc/udev/rules.d/54-gphoto.rules /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/rules.d/54-gphoto.rules 644 123 127
file /etc/udev/rules.d/device-mapper.rules /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/rules.d/device-mapper.rules 644 123 127
file /etc/udev/rules.d/50-firmware.rules /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/rules.d/50-firmware.rules 644 123 127
file /etc/udev/rules.d/75-cd-aliases-generator.rules.optional /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/rules.d/75-cd-aliases-generator.rules.optional 644 123 127
file /etc/udev/rules.d/11-usb-hotplug.rules /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/rules.d/11-usb-hotplug.rules 644 123 127
file /etc/udev/rules.d/97-bluetooth-serial.rules /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/rules.d/97-bluetooth-serial.rules 644 123 127
file /etc/udev/rules.d/75-persistent-net-generator.rules.optional /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/rules.d/75-persistent-net-generator.rules.optional 644 123 127
file /etc/udev/rules.d/60-pcmcia.rules /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/rules.d/60-pcmcia.rules 644 123 127
file /etc/udev/rules.d/99-fuse.rules /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/rules.d/99-fuse.rules 644 123 127
file /etc/udev/rules.d/90-hal.rules /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/udev/rules.d/90-hal.rules 644 123 127
file /etc/fs-version /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/fs-version 644 123 127
file /etc/mtab /data/jenkins/hkrelease_lt00/LT00_SDK_V1.0.3_202411061917/output/lt00/bsp/rootfs/etc/mtab 644 123 127
dir /proc 755 123 127
dir /test 755 123 127
